# DEBIAN/postinst
systemctl enable isitfoggy
systemctl start isitfoggy
